# autostart/conky.desktop
[Desktop Entry]
Name=Conky
TryExec=conky
Exec=conky
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
