# autostart/launchy.desktop
[Desktop Entry]
Name=Launchy
TryExec=launchy
Exec=launchy
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
