# autostart/radiotray.desktop
[Desktop Entry]
Name=radiotray
TryExec=radiotray
Exec=radiotray
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
