# autostart/terminator.desktop
[Desktop Entry]
Name=Terminator
TryExec=terminator
Exec=terminator
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
