# autostart/Whatsapp.desktop
[Desktop Entry]
Name=Whatsapp
TryExec=Whatsapp
Exec=Whatsapp
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
