# autostart/wmname.desktop
[Desktop Entry]
Name=wmname
TryExec=wmname
Exec=wmname LG3D
Terminal=false
Type=Application
StartupNotify=false
OnlyShowIn=awesome
